# tools-pi6ehv/dvbt-rx.service
[Unit]
Description=PI6EHV RB-DVB-T RTL-SDR receiver to SRT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/dvbt-rx
ExecStartPre=/bin/sh -c '! systemctl --user is-active --quiet nicam-rx.service'
ExecStart=%h/dvbt-rx/tools-pi6ehv/dvbt_rx_pi6ehv.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
